# config/kanata/kanata.service
[Unit]
Description=Kanata keyboard remapper
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/bin/kanata --cfg /etc/kanata/kanata.kbd
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
